# meetingbot.conf
description "MeetingBot"
author      "@chrisjaure"

start on (local-filesystems and net-device-up)
stop on shutdown

respawn
respawn limit 20 5

script
  exec /meetingbot/bin/server 2>&1
end script
